# root/etc/cron.daily/ntopng-cleanup
#!/bin/bash

# Delete files not updated in the last 30 days
/usr/bin/find /var/tmp/ntopng/* -mindepth 1 -type f -mtime +30 -delete > /dev/null
